# etc/sv/sshp/run
#!/bin/sh

NAME="sshp"

exec 2>&1

echo "Starting $NAME..."
