# etc/sv/mygoapp/run
#!/bin/sh
echo "starting mygoapp"
exec /opt/mygoapp/mygoapp -c /etc/mygoapp.json  > /dev/null 2>&1
